# myvpn-client/script/if-up.d/20-default-route.sh
#!/bin/bash

IF="$1"
ServerIP="$2"

ip route add $(ip route get $ServerIP | head -n1)
